# build/config/supervisord/named.conf
[program:named]
command=/opt/bind/sbin/named -g -c /opt/bind/etc/named.conf
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
umask=002
priority=2
startretries=3
stopwaitsecs=10
stdout_logfile=/opt/logs/supervisord-named.log
stdout_logfile_maxbytes=0
stderr_logfile=/opt/logs/supervisord-named-error.log
stderr_logfile_maxbytes=0
